# endeavour/etc/rc.local
#!/usr/bin/env bash

sudo kbdrate --delay 300 --rate 50

